# docker-utils/supervisord/mailhog.conf
[program:mailhog]
command=/usr/bin/mailhog
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true
